# rules/90-mtp.rules
SUBSYSTEM=="misc", KERNEL=="mtp_usb", ACTION=="add", RUN+="/usr/bin/chown droidian:droidian /dev/mtp_usb" 
